# deploy/macos/io.speechmesh.apple-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.speechmesh.apple-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/target/release/bridge_tcpd</string>
    <string>--listen</string>
    <string>0.0.0.0:9654</string>
    <string>--bridge-command</string>
    <string>__REPO_ROOT__/bridges/apple-asr/.build/release/apple-asr-bridge</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/SpeechMesh/apple-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/SpeechMesh/apple-bridge.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RUST_LOG</key>
    <string>info</string>
  </dict>
</dict>
</plist>
