# deploy/macos/io.speechmesh.device-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.speechmesh.device-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>__BINARY_PATH__</string>
    <string>agent</string>
    <string>run</string>
    <string>--gateway-url</string>
    <string>__GATEWAY_URL__</string>
    <string>--agent-id</string>
    <string>__AGENT_ID__</string>
    <string>--agent-name</string>
    <string>__AGENT_NAME__</string>
    <string>--device-id</string>
    <string>__DEVICE_ID__</string>
    <string>--provider-id</string>
    <string>device.speaker</string>
    <string>--shared-secret</string>
    <string>__SHARED_SECRET__</string>
    <string>--capability</string>
    <string>speaker</string>
    <string>--reconnect-delay-secs</string>
    <string>2</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/SpeechMesh/device-agent.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/SpeechMesh/device-agent.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RUST_LOG</key>
    <string>info</string>
  </dict>
</dict>
</plist>
